# ansible/ray-head.service.j2
[Unit]
Description=Ray Cluster Service (Head Node)
After=network.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=root
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/ray/venv/bin/ray start --head --port=6379 --dashboard-port=8265 --block
ExecStop=/opt/ray/venv/bin/ray stop
TimeoutStartSec=600
TimeoutStopSec=300
Restart=on-failure
RestartSec=10
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target